# templates/default/etc/init.d/etcd.erb
#!/bin/bash
#
# This shell script takes care of starting and stopping etcd
#
### BEGIN INIT INFO
# Provides:          etcd
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Required-Start:    $network $remote_fs
# Short-Description: Start etcd daemon at boot time
# Description:       Enable etcd service provided by daemon.
### END INIT INFO

# Source function library.
. /lib/lsb/init-functions

# Source etcd configuration
if [ -f /etc/default/etcd ] ; then
 . /etc/default/etcd
fi

NAME="etcd"
PIDFILE=/var/run/$NAME.pid
DAEMON=/usr/local/bin/etcd
LOGDIR=<%= node.etcd.log_dir %>
USER=<%= node.etcd.user %>
GROUP=<%= node.etcd.group %>

start() {
 echo -n "Starting daemon: "$NAME
 start-stop-daemon -C -b -m -p $PIDFILE -c $USER:$GROUP --exec $DAEMON -S > $LOGDIR/etcd 2> $LOGDIR/etcd.err
 echo "."
}

stop() {
 echo -n "Stopping daemon: "$NAME
 start-stop-daemon --retry 10 -K --pidfile $PIDFILE
 echo "."
}

status() {
 start-stop-daemon -T --pidfile $PIDFILE
 if [ $? -eq 0 ]; then
  echo "etcd is running with PID:`cat $PIDFILE`"
 fi
}

reload() {
 # Reload daemons.
 echo -n $"Reloading $prog: "
 kill -1 `pidof etcd`
 RETVAL=$?
 echo
 return $RETVAL
}

# See how we were called.
case "$1" in
 start)
  start
  ;;
 stop)
  stop
  ;;
 restart)
  stop
  start
  ;;
 condrestart)
  stop
  start
  ;;
 status)
  status
  ;;
 reload)
  reload
  ;;
 *)
  echo $"Usage: $0 {start|stop|restart|reload|condrestart|status}"
  exit 1
esac